# deploy/systemd/zentrix-web.service
[Unit]
Description=Zentrix Web API e painel
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/zentrix-web/BackEnd
EnvironmentFile=/opt/zentrix-web/BackEnd/.env
ExecStart=/usr/bin/java -jar /opt/zentrix-web/BackEnd/target/zentrix-web-api-0.1.0-SNAPSHOT.jar
Restart=always
RestartSec=5
SuccessExitStatus=143
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/zentrix-web/BackEnd
CapabilityBoundingSet=
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
